# .config/power-brightness/power-brightness.service
[Unit]
Description=Adjust brightness on power source change

[Service]
Type=oneshot
ExecStart=/usr/local/bin/power-brightness.sh
